# sv-ds03-system/meta-csr-bsp/recipes-core/udev/udev-extraconf/mtpfs.service
[Unit]

Description=mtpfs

[Service]

ExecStart=/etc/udev/scripts/usbmtp.sh
